# nginx-supervisor.ini
[supervisord]
nodaemon=true

[program:nginx]
command = nginx
user = root
autostart = true

[program:php-fpm7]
command = php-fpm7 -F
user = root
autostart = true
